# ansible/resources/crontab
0 */2 * * * /home/ubuntu/store_weather_observations.sh
30 6 * * * /home/ubuntu/store_pollution_data.sh
# end